# client-setup/wmsmon_logrotate.cron
HOME=/
MAILTO=root@localhost

00 2 * * * root /usr/sbin/logrotate -v /etc/wmsmon_logrotate.conf >> /var/log/wmsmon_logrotate.log 2>&1
